# .config/systemd/user/notes-hourly-sync.timer
[Unit]
Description=Notes - sync git repo every hour

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

