# .config/systemd/user/notes-hourly-sync.service
[Unit]
Description=Notes - sync git repo every 3 hours

[Service]
Type=oneshot
ExecStart=%h/nerd-stuff/notes/scripts/sync.fish
