# .config/systemd/user/notes-daily-note.timer
[Unit]
Description=Notes - daily note

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
