# .config/systemd/user/notes-weekly-note.timer
[Unit]
Description=Notes - weekly note

[Timer]
OnCalendar=Sun
Persistent=true

[Install]
WantedBy=timers.target
