# .config/systemd/user/zus-checkin.service
[Unit]
Description=Zus Checkin daily

[Service]
Type=oneshot
ExecStart=%h/Documents/zus_checkin.sh
